# systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=7min
OnActiveSec=1h
Unit=mbsync.service

[Install]
WantedBy=timers.target

